# migration_backup_1753421520/scripts/docker-auto-start.service
[Unit]
Description=Auto-start Docker containers for Elders Guild
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'sg docker -c "docker start elders-guild-postgres-new elders-guild-redis upload-image-service-frontend-1 upload-image-service-backend-1"'
ExecStop=/bin/bash -c 'sg docker -c "docker stop elders-guild-postgres-new elders-guild-redis upload-image-service-frontend-1 upload-image-service-backend-1"'
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
