# config/aicompany-elder.service
[Unit]
Description=AI Company Elder Council Monitoring Service
After=network.target rabbitmq-server.service

[Service]
Type=simple
User=aicompany
Group=aicompany
WorkingDirectory=/home/aicompany/ai_co
Environment="PYTHONPATH=/home/aicompany/ai_co"
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /home/aicompany/ai_co/check_and_fix_workers.py
Restart=always
RestartSec=30
StandardOutput=append:/home/aicompany/ai_co/logs/elder_service.log
StandardError=append:/home/aicompany/ai_co/logs/elder_service_error.log

[Install]
WantedBy=multi-user.target