# configs/aicompany-web.service
[Unit]
Description=AI Company Web Dashboard Service
After=network.target
Wants=network.target

[Service]
Type=forking
User=aicompany
Group=aicompany
WorkingDirectory=/home/aicompany/ai_co
Environment=PYTHONPATH=/home/aicompany/ai_co
Environment=AI_VENV_ACTIVE=1
ExecStartPre=/bin/bash -c 'cd /home/aicompany/ai_co && source venv/bin/activate'
ExecStart=/bin/bash -c 'cd /home/aicompany/ai_co && source venv/bin/activate && python3 web/dashboard_server.py &'
ExecStop=/bin/bash -c 'pkill -f "dashboard_server.py"'
Restart=always
RestartSec=10

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aicompany-web

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/aicompany/ai_co

[Install]
WantedBy=multi-user.target
